# mytv.cron
# run myTV every two hours
0 * * * *     root   /usr/sbin/mytv